# mcp-server/infra/launchd/com.azzas.mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.azzas.mcp</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/docker</string>
        <string>run</string>
        <string>--rm</string>
        <string>--name</string><string>mcp-server</string>
        <string>-p</string><string>3000:3000</string>
        <string>-v</string><string>/Users/artur/bq-analista:/app/repo</string>
        <string>-v</string><string>/etc/mcp:/app/config:ro</string>
        <string>-v</string><string>/var/mcp:/var/mcp</string>
        <string>-e</string><string>MCP_SETTINGS=/app/config/settings.toml</string>
        <string>-e</string><string>MCP_ALLOWLIST=/app/config/allowed_execs.json</string>
        <string>-e</string><string>MCP_BQ_SA_KEY</string>
        <string>-e</string><string>MCP_GITHUB_PAT</string>
        <string>-e</string><string>MCP_AZURE_TENANT_ID</string>
        <string>-e</string><string>MCP_AZURE_CLIENT_ID</string>
        <string>-e</string><string>MCP_AZURE_CLIENT_SECRET</string>
        <string>-e</string><string>MCP_JWT_SECRET</string>
        <string>-e</string><string>MCP_GIT_PUSH=1</string>
        <string>mcp-azzas:latest</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>MCP_BQ_SA_KEY</key><string>__set_via_keychain_wrapper__</string>
        <key>MCP_GITHUB_PAT</key><string>__set_via_keychain_wrapper__</string>
        <key>MCP_AZURE_TENANT_ID</key><string>__set_via_keychain_wrapper__</string>
        <key>MCP_AZURE_CLIENT_ID</key><string>__set_via_keychain_wrapper__</string>
        <key>MCP_AZURE_CLIENT_SECRET</key><string>__set_via_keychain_wrapper__</string>
        <key>MCP_JWT_SECRET</key><string>__set_via_keychain_wrapper__</string>
    </dict>
    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><true/>
    <key>StandardOutPath</key><string>/var/log/mcp/stdout.log</string>
    <key>StandardErrorPath</key><string>/var/log/mcp/stderr.log</string>
</dict>
</plist>
