# mcp-server/infra/launchd/com.azzas.mcp.alerts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.azzas.mcp.alerts</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/docker</string><string>exec</string>
        <string>mcp-server</string>
        <string>python</string><string>-m</string><string>mcp_exec.alerts</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict><key>Minute</key><integer>0</integer></dict>
    <key>StandardOutPath</key><string>/var/log/mcp/alerts.log</string>
    <key>StandardErrorPath</key><string>/var/log/mcp/alerts.err</string>
</dict>
</plist>
